# deploy/supervisor-queue-worker.conf
; Supervisor config for Rabie Co's database-driven queue worker.
; Processes SendOrderNotificationJob / SendContactNotificationJob (and any
; future queued jobs) so admin emails never block the checkout/contact
; request, and failed sends retry with backoff instead of being lost.
;
; Install:
;   sudo cp deploy/supervisor-queue-worker.conf /etc/supervisor/conf.d/rabie-queue-worker.conf
;   # edit command= below to match the real deploy path and PHP binary first
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start rabie-queue-worker:*
;
; Inspect:
;   sudo supervisorctl status rabie-queue-worker:*
;   tail -f /var/log/rabie-queue-worker.log
;
; Failed jobs land in the `failed_jobs` table after exhausting retries:
;   php artisan queue:failed
;   php artisan queue:retry <id>|all

[program:rabie-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/lampp/bin/php /var/www/rabie-co/artisan queue:work database --sleep=3 --tries=5 --max-time=3600
directory=/var/www/rabie-co
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/rabie-queue-worker.log
stopwaitsecs=3600
